# install_os/etc/inittab
id:5:initdefault:

ca::ctrlaltdel:/sbin/reboot -f

l0:0:wait:/sbin/halt -f
l6:6:wait:/etc/reboot -f

si:S:sysinit:/etc/boot-up

c2:2345:respawn:/sbin/sh < /dev/tty2 > /dev/tty2 2>/dev/tty2
c3:2345:respawn:/sbin/sh < /dev/tty3 > /dev/tty3 2>/dev/tty3

rd:5:respawn:/usr/bin/run_dialog < /dev/tty5 > /dev/tty5 2>/dev/tty5
