# baseos/rootfs/common/ram/inittab
id:2:initdefault:

si:S:sysinit:/etc/rc.d/boot
